# unit.service
[Unit]
Description=Jobboard FastAPI (uvicorn)
After=network.target

[Service]
Type=simple
User=jobsvc
Group=jobsvc
WorkingDirectory=/opt/jobboard
# Экспортируем окружение (опционально). Замени DB URL если нужно.
Environment="PATH=/opt/jobboard/venv/bin"
Environment="DATABASE_URL=postgresql+psycopg2://jobuser:jobpass@localhost:5432/jobboard"
ExecStart=/opt/jobboard/venv/bin/uvicorn app.main:app --host 0.0.0.0 --port 8000 --workers 1
Restart=always
RestartSec=5
LimitNOFILE=65536

[Install]
WantedBy=multi-user.target
