# unit.service
[Unit]
Description=Jobboard FastAPI
After=network.target

[Service]
User=jobsvc
Group=jobsvc
WorkingDirectory=/path/to/jobboard
Environment=DATABASE_URL=postgresql+psycopg2://jobuser:jobpass@localhost:5432/jobboard
ExecStart=/path/to/venv/bin/uvicorn app.main:app --host 0.0.0.0 --port 8000
Restart=always

[Install]
WantedBy=multi-user.target