# units/isolated_jobboard.service
[Unit]
Description=Jobboard FastAPI (uvicorn)
After=network.target

[Service]
Type=simple
User=jobsvc
Group=jobsvc
WorkingDirectory=/opt/jobboard

EnvironmentFile=/opt/jobboard/.env
Environment="PATH=/opt/jobboard/venv/bin"

ExecStart=/opt/jobboard/venv/bin/uvicorn app.main:app --host 127.0.0.1 --port 8000 --workers 1
Restart=always
RestartSec=5
LimitNOFILE=65536

[Install]
WantedBy=multi-user.target
